# files/easectl.ini
[program:easelm]
command=/opt/ease/bin/easelm
stdout_logfile=/opt/ease/log/easelm.log
autostart=true
autorestart=true
startsecs=3
priority=1
stopasgroup=true
killasgroup=true

[program:easeserver]
command=/opt/ease/bin/easeserver /opt/ease/etc/easeserver.ini
stdout_logfile=/opt/ease/log/easeserver.log
autostart=true
autorestart=true
startsecs=3
priority=1
stopasgroup=true
killasgroup=true

[program:easeagent]
command=/opt/ease/bin/easeagent /opt/ease/etc/easeagent.ini
stdout_logfile=/opt/ease/log/easeagent.log
autostart=true
autorestart=true
startsecs=3
priority=1
stopasgroup=true
killasgroup=true

[program:easeweb]
command=/opt/ease/bin/easeweb -w /opt/ease/www -c /opt/ease/etc/easeweb.json
stdout_logfile=/opt/ease/log/easeweb.log
autostart=true
autorestart=true
startsecs=3
priority=1
stopasgroup=true
killasgroup=true

[program:easelog]
command=/opt/ease/bin/easelog /opt/ease/etc/easelog.ini
stdout_logfile=/opt/ease/log/easelog.log
autostart=true
autorestart=true
startsecs=3
priority=1
stopasgroup=true
killasgroup=true
